# agent-ovs/debian/opflex-agent.upstart
description "Opflex Agent"

start on (local-filesystems)
stop on runlevel [!2345]

kill signal SIGINT

respawn
respawn limit 10 5

pre-start script
    test -x /usr/bin/opflex-agent || { stop; exit 0; }
    test -e /etc/opflex-agent-ovs/opflex-agent-ovs.conf || { stop; exit 0; }
end script

exec /usr/bin/agent_ovs --syslog \
     -c /etc/opflex-agent-ovs/opflex-agent-ovs.conf \
     -c /etc/opflex-agent-ovs/plugins.conf.d \
     -c /etc/opflex-agent-ovs/conf.d
